# supervisor.example.conf
[program:athena]
environment=MYSQL_PASSWORD="******"
directory=/path/to/project/
command=make start-gunicorn
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/athena.log

